# ansible_stuffs/hands_on/k8s_operations/k8s-cluster-creation_with_ansible_module/roles/etcd/templates/etcd.service.j2
[Unit]
Description=etcd
Documentation=https://github.com/etcd-io/etcd
After=network.target

[Service]
Type=notify
ExecStart=/usr/local/bin/etcd --config-file=/etc/etcd/etcd.conf.yaml
Restart=on-failure
RestartSec=5
LimitNOFILE=1048576
LimitNPROC=1048576
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target